# packages/wirikiki/wirikiki-1.1.2-py3-none-any.whl/wirikiki/service/wirikiki.service
[Unit]
Description=Wirikiki Wiki

[Service]
Type=simple
Environment="FG=1"
WorkingDirectory=/home/fab/Documents/mywiki/
ExecStart=wirikiki

[Install]
WantedBy=default.target
